# data_crawling/crontab
0 0 * * * cd /clear-closet && /usr/local/bin/python3 crawling_main.py >> /clear-closet/cron.log 2>&1
